# 9/postfix-supervisord.conf
[program:postfix]
process_name = master
directory = /etc/postfix
command =  /etc/init.d/postfix start
startsecs = 0
autorestart = false
