# 9/syslog-ng-supervisor.conf
[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground --no-caps
stdout_syslog=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
#stdout_capture_maxbytes=1MB
