# gvm20/rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv sh
redis-server /etc/redis/redis.conf