# gvm20/rootfs/etc/services.d/gsad/run
#!/usr/bin/with-contenv sh

# s6-setuidgid gvm gsad /opt/gvm/sbin/gsad