# gvm20/rootfs/etc/services.d/gsad/run
#!/usr/bin/with-contenv sh

s6-setuidgid gvm /opt/gvm/sbin/gsad